# supervisord_dev.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:celery]
command=celery -A wuvt.trackman.tasks.celery worker -B -l info

[program:uwsgi]
command=uwsgi --ini uwsgi.ini:dev --die-on-term
