# supervisord_dev.conf
[supervisord]
nodaemon=true

[program:celery]
command=celery -A wuvt.trackman.tasks.celery worker -B -l info

[program:uwsgi]
command=uwsgi --yaml uwsgi.yml:dev --die-on-term
